# systemd/quadlets/corporatetraveldc-ops-brief.container
# One-shot container 芒聙聰 run by systemd timer at 00:00, 06:00, 12:00, 18:00 ET
[Unit]
Description=Corporate Travel DC 芒聙聰 6-hour ops brief (00/06/12/18 ET)

[Container]
Image=localhost/corporatetraveldc-poller:latest
Volume=/var/lib/corporatetraveldc:/var/lib/corporatetraveldc:z
EnvironmentFile=/etc/corporatetraveldc/dispatch.env
EnvironmentFile=/etc/corporatetraveldc/dispatch-secrets.env
Exec=python3 poller/skills/ops_brief.py

[Service]
Type=oneshot
TimeoutStartSec=180
