# systemd/quadlets/corporatetraveldc-daily-opsplan.container
[Unit]
Description=Corporate Travel DC 芒聙聰 daily ATCSCC ops plan snapshot (07:00 ET)

[Container]
Image=localhost/corporatetraveldc-poller:latest
Volume=/var/lib/corporatetraveldc:/var/lib/corporatetraveldc:z
EnvironmentFile=/etc/corporatetraveldc/dispatch.env
EnvironmentFile=/etc/corporatetraveldc/dispatch-secrets.env
Exec=python3 poller/fetchers/atcscc_opsplan.py

[Service]
Type=oneshot
TimeoutStartSec=120
