# systemd/quadlets/corporatetraveldc-ultrafeeder.container
[Unit]
Description=Corporate Travel DC - UltraFeeder ADS-B (1090 MHz)
After=network-online.target

[Container]
Image=ghcr.io/sdr-enthusiasts/docker-adsb-ultrafeeder:latest
ContainerName=corporatetraveldc-ultrafeeder
AddDevice=/dev/rtl_sdr_adsb:/dev/rtl_sdr_adsb
Volume=/dev/bus/usb:/dev/bus/usb:ro
Volume=/var/lib/corporatetraveldc:/var/lib/corporatetraveldc:z
EnvironmentFile=/etc/corporatetraveldc/dispatch.env
EnvironmentFile=/etc/corporatetraveldc/dispatch-secrets.env
Environment=READSB_DEVICE_TYPE=rtlsdr
Environment=READSB_RTLSDR_DEVICE=ADSB1090
Environment=READSB_LAT=38.87485
Environment=READSB_LON=-77.14410
Environment=READSB_ALT=85m
Environment=READSB_MODEAC=true
Environment=READSB_MLAT=true
Environment=READSB_NET_ENABLE=true
Environment=READSB_NET_ONLY=false
Environment=TAR1090_DEFAULTCENTERLAT=38.87485
Environment=TAR1090_DEFAULTCENTERLON=-77.14410
Environment=TAR1090_DEFAULTRANGE=80
Environment=TAR1090_PAGETITLE=corporatetraveldc-kdca
Environment=ULTRAFEEDER_CONFIG=adsb,piaware.flightaware.com,30005,beast_reduce_plus_out,uuid=807b07ce-059d-457f-ba8b-d75725679565;adsb,data-out.flightradar24.com,10001,beast_reduce_out;adsb,data.adsbhub.org,5002,beast_reduce_out;mlat,collector.opensky-network.org,30003
PublishPort=8080:80
PublishPort=30005:30005
PublishPort=127.0.0.1:30003:30003
Environment=UUID=784fdeaa-57f1-42bc-936a-cf562775e8bd
Environment=OPENSKY_USERNAME=corporatetraveldc
AutoUpdate=registry

[Service]
Restart=always
RestartSec=30
TimeoutStartSec=60

[Install]
WantedBy=default.target
