# cloudflared/cloudflared.service.d/restart.conf
[Service]
Restart=always
RestartSec=10
StartLimitIntervalSec=0
