# systemd/unbound-anchor.service.d/reload.conf
# Drop-in for unbound-anchor.service (the daily timer target)
# Place at: /etc/systemd/system/unbound-anchor.service.d/reload.conf
#
# Replaces the upstream ExecStart with our wrapper script, which only
# signals Unbound to reload if the anchor actually changed (exit code 1).
# This ensures Unbound is NOT disturbed on days when the anchor is current.

[Service]
# Clear the upstream ExecStart, then set ours
ExecStart=
ExecStart=/usr/local/sbin/unbound-anchor-update
# Both exit codes are success (0=current, 1=updated — wrapper normalizes to 0)
SuccessExitStatus=0
