# .config/containers/systemd/corporatetraveldc-dumpvdl2.container
[Unit]
Description=Corporate Travel DC - dumpvdl2 VDL2 decoder
After=network-online.target
After=corporatetraveldc-acarsrouter.service
Wants=corporatetraveldc-acarsrouter.service

[Container]
Image=ghcr.io/sdr-enthusiasts/docker-dumpvdl2:latest
ContainerName=corporatetraveldc-dumpvdl2
AddDevice=/dev/rtl_sdr_acars:/dev/rtl_sdr_acars
Volume=/dev/bus/usb:/dev/bus/usb:ro
EnvironmentFile=/etc/corporatetraveldc/dispatch.env
EnvironmentFile=/etc/corporatetraveldc/dispatch-secrets.env
Environment=TZ=America/New_York
Environment=RTL_SERIAL=ACARS0130
Environment=FEED_ID=corporatetraveldc-kdca-vdl2
Environment=FREQUENCIES=136.650;136.800;136.900;136.925;136.975
Environment=GAIN=40
Environment=ZMQ_MODE=server
Environment=ZMQ_ENDPOINT=tcp://0.0.0.0:45555
Environment=VDLM_FILTER_ENABLE=TRUE
Network=acars-net
AutoUpdate=registry

[Service]
Restart=always
RestartSec=30
TimeoutStartSec=60

[Install]
WantedBy=default.target
