# .config/containers/systemd/ntfy.container
[Unit]
Description=ntfy — primary push broker (container, port 2586)
After=network-online.target
StartLimitBurst=3
StartLimitIntervalSec=120s
OnFailure=ntfy-container-alert.service

[Container]
Image=docker.io/binwiederhier/ntfy:v2.25.0
ContainerName=ntfy
PublishPort=2586:2586
Volume=/etc/ntfy:/etc/ntfy:ro,z
Volume=/var/lib/ntfy-backup:/var/lib/ntfy:z

# SMTP + Twilio credentials (from dispatch-secrets.env)
EnvironmentFile=/etc/corporatetraveldc/dispatch-secrets.env
Exec=serve

[Service]
Restart=on-failure
RestartSec=30s

[Install]
WantedBy=default.target
